# service/watch/run
#!/bin/sh

cd /watch

exec 2>&1 /usr/bin/python ./watch.py
